# misc./etc-systemd-system/system/yay-update.service
[Unit]
Description=Update AUR Package Database

[Service]
Type=oneshot
ExecStart=/usr/bin/yay -Sy --noconfirm
